# crontab
0 9 * * * sudo systemctl restart media-autoplay
